# scripts/nl.botlease.newsbot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- BotLease news bot — weekly LaunchAgent.
     INSTALL (you run this once):
       cp scripts/nl.botlease.newsbot.plist ~/Library/LaunchAgents/
       launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/nl.botlease.newsbot.plist
     TEST a run now:
       launchctl kickstart -p gui/$(id -u)/nl.botlease.newsbot
     UNINSTALL:
       launchctl bootout gui/$(id -u)/nl.botlease.newsbot
       rm ~/Library/LaunchAgents/nl.botlease.newsbot.plist
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>nl.botlease.newsbot</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/werk/Documents/Python/botlease.nl/scripts/run_news_weekly.sh</string>
    </array>

    <!-- Weekly: Monday 09:00 local time -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <!-- launchd has a minimal env; pin PATH and HOME. -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/werk</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>/Users/werk/Documents/Python/botlease.nl</string>

    <key>StandardOutPath</key>
    <string>/Users/werk/Documents/Python/botlease.nl/logs/launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/werk/Documents/Python/botlease.nl/logs/launchd.err.log</string>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
